# irys-rlm.service
[Unit]
Description=Irys RLM API Service
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/legal-rlm
Environment="PATH=/home/ubuntu/legal-rlm/.venv/bin"
ExecStart=/home/ubuntu/legal-rlm/.venv/bin/python run_server.py
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=irys-rlm

[Install]
WantedBy=multi-user.target
